# deploy/systemd/nerubot.service
[Unit]
Description=NeruBot Discord Music Bot
After=network.target
Wants=network.target

[Service]
Type=simple
User=nerubot
Group=nerubot
WorkingDirectory=/home/nerubot/nerubot
Environment=PYTHONPATH=/home/nerubot/nerubot
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/nerubot/nerubot/nerubot_env/bin/python src/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nerubot

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/nerubot/nerubot /home/nerubot/logs
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
